# ansible/roles/iptables/files/emergency-iptables-recovery.service
[Unit]
Description=Emergency iptables recovery (iptables-persistent stop + permissive chains)
Documentation=man:iptables(8)
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/emergency-iptables-recovery.sh


